# crontab
MAILTO=""
*/1 * * * * /bin/countcommit.sh
